# incoming/bhphoto-category-discovery/ERIP_BHPhoto_Retail_Ingest_GreenTier1_v1/systemd/erip-bhphoto-ingest.timer
[Unit]
Description=Daily ERIP B&H Photo Retail Ingest
[Timer]
OnCalendar=*-*-* 03:35:00
Persistent=true
RandomizedDelaySec=300
AccuracySec=1min
Unit=erip-bhphoto-ingest.service
[Install]
WantedBy=timers.target
